# gps-bridge-connect.service
[Unit]
Description=gps-bridge GPS relay connector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Token is stored in ~/.gps-bridge/config.json (saved on first `gps-bridge connect --token <TOKEN>`)
# No EnvironmentFile or connect.env needed.
# Replace the path below with the output of: which gps-bridge
ExecStart=/usr/local/bin/gps-bridge connect
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
